# ubuntu-installer/etc/udev/rules.d/95-monitor-hotplug.rules
KERNEL=="card0", SUBSYSTEM=="drm", ENV{DISPLAY}=":0", ENV{XAUTHORITY}="/home/alex/.Xauthority", RUN+="/home/alex/.config/i3/screen-connected"
